# .local/bin/prevent-sleep.service
[Unit]
Description=Prevent system from sleeping when lid is closed

[Service]
ExecStart=/home/tetr-a/.local/bin/prevent-sleep.sh

[Install]
WantedBy=multi-user.target
